# code/set-up-right-click/etc-udev-rules.d/70-touchscreen-waveshare.rules
SUBSYSTEMS=="usb",ACTION=="add",KERNEL=="event*",ATTRS{idVendor}=="0eef",ATTRS{idProduct}=="0005",SYMLINK+="twofingtouch",RUN+="/bin/chmod a+r /dev/twofingtouch"
KERNEL=="event*",ATTRS{name}=="WaveShare WS170120",SYMLINK+="twofingtouch",RUN+="/bin/chmod a+r /dev/twofingtouch"	